# squid.ini
[program:squid]
command=/usr/sbin/squid -N -d -f /config/squid/squid.conf
