# squid.ini
[program:squid]
command=/usr/sbin/squid -N -d -f /etc/squid/squid.conf
